# deploy/systemd/polymarket-sampler.service
[Unit]
Description=Polymarket continuous paper sampling
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-wallet-tracker
EnvironmentFile=/etc/polymarket-wallet-tracker.env
ExecStart=/opt/polymarket-wallet-tracker/.venv/bin/python -m backend.scripts.run_continuous_sampling --repeat-seconds 60 --max-cycles 0 --wallet-limit 50 --trade-page-limit 100 --trade-max-pages 2 --token-limit 30 --paper-lookback-minutes 120
Restart=always
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
