# deploy/systemd/polymarket-maintenance.timer
[Unit]
Description=Run Polymarket sampling maintenance every six hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
