# deploy/systemd/polymarket-backup.timer
[Unit]
Description=Back up Polymarket PostgreSQL database daily

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
